# cicd/packages/fmlx-spi2can/DEBIAN/postinst
# Start systemd service
sudo systemctl start fmlx-spi2can
sudo systemctl enable fmlx-spi2can