# cicd/packages/fmlx-spi2can/DEBIAN/prerm
sudo systemctl stop fmlx-spi2can
sudo systemctl disable fmlx-spi2can
sudo systemctl daemon-reload